# home/dot_config/autostart/autocutsel.desktop
[Desktop Entry]
Type=Application
Name=Clipboard to Primary Sync
Comment=One-way sync from CLIPBOARD to PRIMARY selection
Exec=sh -c "while true; do xsel -bo 2>/dev/null | xsel -pi 2>/dev/null; sleep 0.5; done"
Terminal=false
Hidden=false
X-GNOME-Autostart-enabled=true
