# home/dot_config/systemd/user/cleanup-orphaned-procs.timer
[Unit]
Description=Periodically clean up orphaned TUI processes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
